# bnx2x/untangle/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

iptables -I INPUT -p igmp -j ACCEPT 
iptables -I FORWARD -i eth0.36 -s 10.2.0.0/16 -d 239.0.0.0/8 -p udp -j ACCEPT

exit 0
